# packaging/debian/telephone-booth-tailscale-serve.service
[Unit]
Description=Tailscale serve proxy for Telephone Booth debug surface
After=tailscaled.service telephone-booth.service
Requires=tailscaled.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/tailscale serve --bg --https=443 --set-path=/ http://127.0.0.1:8080
ExecStop=/usr/bin/tailscale serve reset

[Install]
WantedBy=multi-user.target
