# packaging/debian/telephone-booth.service
[Unit]
Description=Telephone Booth phone client runtime
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=notify
User=phonebooth
Group=audio
# `phonebooth` is required in addition to the runtime device groups: overriding
# the primary group with Group=audio drops the user's own group, which owns the
# 0750 /etc/phone-booth config directory. Without it the service cannot read
# config.toml; the `check` pre-start step (ExecStartPre) then fails, so the unit
# refuses to start instead of silently running on built-in defaults.
SupplementaryGroups=gpio dialout phonebooth
ExecStartPre=/usr/bin/telephone-booth check
ExecStart=/usr/bin/telephone-booth run
Restart=on-failure
RestartSec=2s
WatchdogSec=30s
EnvironmentFile=/etc/phone-booth/env
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/phone-booth /var/log/phone-booth
CapabilityBoundingSet=
NoNewPrivileges=true
RestrictRealtime=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
